# ansible/roles/backup_timer/templates/server-apps-backup.service.j2
[Unit]
Description=Server Apps Backup Service
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory={{ server_apps_project_dir }}
ExecStart={{ server_apps_project_dir }}/.venv/bin/python runbook/backup.py
User={{ backup_service_user }}
Group=docker
Environment=PYTHONUNBUFFERED=1
Environment=RUNBOOK_LOCK_TIMEOUT={{ backup_lock_timeout | default(300) }}
StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-apps-backup

[Install]
WantedBy=multi-user.target
